# tools/dingtalk-requirement-inbox/deploy/trinity-requirement-inbox.service
[Unit]
Description=Trinity DingTalk Requirement Inbox (Stream)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trinity
Group=trinity
WorkingDirectory=/opt/trinity-tool/dingtalk-requirement-inbox
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/opt/trinity-tool/dingtalk-requirement-inbox/.env
ExecStart=/opt/trinity-tool/dingtalk-requirement-inbox/.venv/bin/python -m src.main stream
Restart=always
RestartSec=10
# Dedup SQLite under data/
ReadWritePaths=/opt/trinity-tool/dingtalk-requirement-inbox/data

[Install]
WantedBy=multi-user.target
